# supervisord.conf
[program:checkalerts]
command = /home/dotcloud/current/hello_world.php
