# jobs/crontab
*/1 * * * * "~/current/jobs/checkAlert.php" >> test.txt